# disordered/sh/startup/supervisord.d/robot.ini
[group:peak]
programs=roscore,test_node

[program:roscore]
command=/home/peak/robot/develop/ros_test/src/test_node/startup/roscore.sh
autostart=true
user=peak
redirect_stderr=true
stdout_logfile=/home/peak/.ros/roscore.log
loglevel=info
stopsignal=TERM
stopasgroup=true

[program:test_node]
command=/home/peak/robot/develop/ros_test/src/test_node/startup/test_node.sh
autostart=true
user=peak
redirect_stderr=true
stdout_logfile=/home/peak/.ros/test_node.log
loglevel=info
stopsignal=TERM
stopasgroup=true
